# .online/crond-supervisor.conf
[program:crond]
command=crond -f
autorestart=true
autostart=true
priority=100
