# .online/crond-supervisor.conf
[program:crond]
command=crond -f -d 6
autorestart=true
autostart=true
priority=100
